# kinit-refresh.timer
[Unit]
Description=Refresh a kerberos ticket every hour

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
