# kinit-fedora.timer
[Unit]
Description=Refresh a kerberos ticket every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=6days

[Install]
WantedBy=timers.target
